# debian/postinst
#!/bin/bash
set -e

# Create vectorizer user and group
if ! getent group vectorizer >/dev/null; then
    groupadd --system vectorizer
fi

if ! getent passwd vectorizer >/dev/null; then
    useradd --system --gid vectorizer --home-dir /var/lib/vectorizer \
        --shell /usr/sbin/nologin --comment "Vectorizer service user" vectorizer
fi

# Create directories
mkdir -p /var/lib/vectorizer/data
mkdir -p /var/lib/vectorizer/snapshots
mkdir -p /var/log/vectorizer

# Set ownership
chown -R vectorizer:vectorizer /var/lib/vectorizer
chown -R vectorizer:vectorizer /var/log/vectorizer
chmod 755 /var/lib/vectorizer
chmod 755 /var/lib/vectorizer/data
chmod 755 /var/lib/vectorizer/snapshots
chmod 755 /var/log/vectorizer

# Set config permissions
if [ -f /etc/vectorizer/config.yml ]; then
    chmod 644 /etc/vectorizer/config.yml
fi

echo "Vectorizer installed successfully!"
echo "Start the service with: sudo systemctl start vectorizer"
echo "Enable at boot with: sudo systemctl enable vectorizer"

#DEBHELPER#

